# production/minfee.crontab
@reboot /usr/local/bin/litecoind >/dev/null 2>&1
