# production/minfee.crontab
@reboot /usr/local/bin/ferrited >/dev/null 2>&1
